# data/varlink-http-bridge.service.in
[Unit]
Description=Varlink HTTP Bridge
After=network.target

[Service]
Type=simple
ExecStart=@bindir@/varlink-http-bridge --bind 0.0.0.0:1031
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
